# files/etc/uci-defaults/10_system_setup
#!/bin/sh
uci set system.@system[0].timezone='MSK-3'
uci set system.@system[0].zonename='Europe/Moscow'
uci commit system

HASH=$(openssl passwd -6 '20057771925jk')
passwd -p "$HASH" root 2>/dev/null

exit 0
