# files/etc/uci-defaults/70_enable_zapret
#!/bin/sh
/etc/init.d/zapret enable
exit 0
