# files/etc/uci-defaults/71_zapret_clear_hosts
#!/bin/sh
: > /opt/zapret/ipset/zapret-hosts-user.txt
: > /opt/zapret/ipset/zapret-hosts-user-exclude.txt
exit 0
